# supervisord.conf
[supervisord]
nodaemon=true

[program:duoauthproxy]
command=/opt/duoauthproxy/bin/authproxyctl start
